# .deploy/config/crontab
* * * * * php /var/www/html/panel/storedcard.php